# rcLocal/nod1/rc.local
#!/bin/sh
ucarp -i ens33 -s 192.168.13.1 -v 1 -p secret -a 192.168.13.254 -u /etc/vip-up.sh -d /etc/vip-down.sh -z -B

iptables -t nat --flush

iptables -t nat -A POSTROUTING -o ens34 -s 192.168.13.0/24 -j SNAT --to-source 172.16.1.13

conntrack -F
exit 0
